# services/autorestic@laptop.timer
[Timer]
# Wait for 30 seconds after log in
OnStartupSec=30s

[Install]
WantedBy=timers.target