# services/bazarr.timer
[Timer]
# Wait for 5 seconds after boot
OnStartupSec=5s

[Install]
WantedBy=timers.target